# dist/pipydailyclock.service
[Unit]
Description=PiPyDailyClock
After=multi-user.target

[Service]
User=pi
Group=pi
Type=idle
WorkingDirectory=/home/pi/PiPyDailyClock
ExecStart=/home/pi/PiPyDailyClock/venv/bin/python3 pipydailyclock.py
Restart=on-failure
RestartSec=5
Nice=1

[Install]
WantedBy=multi-user.target
